# roles/common/files/node_exporter.service
[Unit]
Description=NodeExporter
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/node_exporter --collector.systemd --collector.processes --web.listen-address=":59100"
Restart=on-failure

[Install]
WantedBy=multi-user.target